# conf/supervisord.conf
[supervisord]
childlogdir = /var/tmp/logs/supervisor/
logfile = /var/tmp/logs/supervisor/supervisord.log
pidfile = /tmp/supervisord-scoringservice.pid
logfile_maxbytes = 50MB
logfile_backups=100
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
nocleanup = true
strip_ansi = false
[unix_http_server]

file=/tmp/supervisor.sock

chmod=0777
# environment = KEY1="value1",KEY2="value2"

# section must remain in the configuration for the standard setup of supervisor to work properly
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface



#A TCP host:port value or on which supervisor will listen for HTTP/XML-RPC requests
[inet_http_server]
port = 127.0.0.1:9001

[supervisorctl]
serverurl = http://localhost:9001


[program:url_collection_worker]
directory=%(here)s/../../
command=celery --app=pipeline_url worker --loglevel=info -c 40 -Q data_collection_queue
autostart=false
autorestart=true
stdout_logfile=/var/tmp/logs/scraping/celery/workers-stdout.log
stderr_logfile=/var/tmp/logs/scraping/celery/workers-stdout.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB